# scripts/etc/init.d/rpi3-router-dhcpd
#!/sbin/openrc-run

container_name="dhcpd"

depend() {
        need rpi3-router-base
}

get_containerid() {
	containerid=$(docker ps -aqf "name=$container_name")
}

start() {
	ebegin "Starting dhcpd container"
	docker run -d --name $container_name --mount source=rpi3-dhcpd-data,target=/data --network=lan rpi3-dhcpd
        if [[ $? != 0 ]]; then
		eend "Failed to start container. Please remove all $container_name container manually"
		exit 1
	fi
}

stop() {
	ebegin "Stopping dhcpd container"
	get_containerid
        if [ ! -z $containerid ]; then
		docker stop $containerid
		docker rm $containerid
	fi
}
