# scripts/etc/init.d/rpi3-router-passthrough-wan
#!/sbin/openrc-run

container_name="passthrough-wan"

depend() {
        need rpi3-router-base
}

get_container_id() {
	container_id=$(docker ps -aqf "name=$container_name")
}

start() {
	ebegin "Starting $container_name container"
	docker create --name $container_name --cap-add=NET_ADMIN --network=lantowan1 --ip=10.1.1.2 rpi3-passthrough-wan
        if [[ $? != 0 ]]; then
		eend "Failed to start container. Please remove all $container_name container manually"
		exit 1
	fi
	get_container_id
	docker network connect wan --ip=10.1.0.1 $container_id
	docker start $container_id
}

stop() {
	ebegin "Stopping $container_name container"
	get_container_id
        if [ ! -z $container_id ]; then
		docker stop $container_id
		docker rm $container_id
	fi
}
