# scripts/etc/init.d/rpi3-router-core
#!/sbin/openrc-run

container_name="router-core"

depend() {
        need rpi3-router-base
}

get_container_id() {
	container_id=$(docker ps -aqf "name=$container_name")
}

start() {
	ebegin "Starting $container_name container"
	docker create --name $container_name --cap-add=NET_ADMIN --network lan --ip 10.1.1.10 rpi3-router
        if [[ $? != 0 ]]; then
		eend "Failed to start container. Please remove all $container_name container manually"
		exit 1
	fi
	get_container_id
	docker network connect lantowan0 $container_id
	docker network connect lantowan1 $container_id
	docker start $container_id
}

stop() {
	ebegin "Stopping $container_name container"
	get_container_id
        if [ ! -z $container_id ]; then
		docker stop $container_id
		docker rm $container_id
	fi
}
