# scripts/etc/init.d/rpi3-router-base
#!/sbin/openrc-run

start() {
	ebegin "Setting up rpi3-router networking"
        # Had trouble with this script being executed before docker daemon startup was completed
	sleep 10
        # Dont send or accept redirects
	sysctl -w net.ipv4.conf.all.accept_redirects=0
	sysctl -w net.ipv4.conf.all.send_redirects=0
	sysctl -w net.ipv6.conf.all.accept_redirects=0
	# sysctl -w net.ipv6.conf.all.send_redirects=0
	# enable ipv6
        sysctl -w net.ipv6.conf.all.disable_ipv6=0
        sysctl -w net.ipv6.conf.default.disable_ipv6=0
	# Not sure atm why rp filter needs to be turned off. Not doing any ip spoofing here.
	sysctl -w net.ipv4.conf.all.rp_filter=0
	sysctl -w net.ipv4.conf.default.rp_filter=0
	sysctl -w net.ipv4.conf.lo.rp_filter=0
        # Fix bridge mac address otherwise windows detects a "new" network with every restart
	ip link set lanbridge address 02:42:fd:17:01:0d
        # This bridge is connected to another bridge in hostapd so we rather enable the spanning tree protocol
        brctl stp lanbridge yes
        # Remove all stopped container
	docker rm $(docker ps -a -q)

	iptables -I DOCKER-ISOLATION-STAGE-1 -i lanbridge -o lanbridge -j ACCEPT
        iptables -I DOCKER-ISOLATION-STAGE-1 -i lantowan0bridge -o lantowan0bridge -j ACCEPT
        iptables -I DOCKER-ISOLATION-STAGE-1 -i lantowan1bridge -o lantowan1bridge -j ACCEPT
}

